# systemd/prepperpi-monitor.service
[Unit]
Description=PrepperPi System Monitor
After=network.target

[Service]
Type=simple
User=prepperpi
Group=prepperpi
ExecStart=/opt/prepperpi/venv/bin/python3 /opt/prepperpi/scripts/monitor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
